# cookbooks/jboss/templates/jboss.erb
#!/bin/sh  
#  
# JBoss standalone control script  
#  
# chkconfig: 234 80 20  
# description: JBoss AS Standalone  
# processname: standalone  
# pidfile: /var/run/jboss-as/jboss-as-standalone.pid  
# config: /etc/jboss-as/jboss-as.conf  
 
 # Source function library.  
. /etc/init.d/functions  
  
# Load Java configuration.  
[ -r /etc/java/java.conf ] && . /etc/java/java.conf  
export JAVA_HOME  
   
JBOSS_USER=jboss  
export JBOSS_USER  
     
# Load JBoss AS init.d configuration.  
if [ -z "$JBOSS_CONF" ]; then  
  JBOSS_CONF="/etc/jboss-as/jboss-as.conf"  
fi  
         
[ -r "$JBOSS_CONF" ] && . "${JBOSS_CONF}"  
          
# Set defaults.  
            
if [ -z "$JBOSS_HOME" ]; then  
JBOSS_HOME=/opt/jboss  
fi  
export JBOSS_HOME




